# install/global/tunnel.conf
[program:tunnel]
command=/root/create_tun.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/tunnel.err.log
stdout_logfile=/var/log/tunnel.out.log
